# ndn-server/conf/etc/service/sidecar/run
#!/bin/sh

exec sidecar -public ${IPADDR}
