# fileserve/conf/etc/service/fileserve/run
#!/bin/sh

cd /opt/fileserve
exec python3 ./app.py 2>&1
