# deploy/traffic/conf/service/ndn-traffic-client/run
#!/bin/sh

sleep 2
exec ndn-traffic-client /usr/local/etc/ndn/ndn-traffic-client.conf
