# deploy/nfd/conf/etc/service/ndn-autoconfig/run
#!/bin/sh

sleep 2
exec ndn-autoconfig 2>&1
