# demo/secondary/conf/etc/service/secondary/run
#!/bin/sh

exec secondary -primary $PRIMARY 2>&1
