# deploy/traffic/conf/etc/service/ndn-traffic-server/run
#!/bin/sh

sleep 2
exec ndn-traffic-server -q /usr/local/etc/ndn/ndn-traffic-server.conf
