# deploy/primary/conf/etc/service/primary/run
#!/bin/sh

exec primary -watcher $WATCHER 2>&1
