# deploy/primary/conf/etc/service/primary/run
#!/bin/sh

exec primary -watcher $WATCHER -addr $PUBLIC_ADDR 2>&1
