# ndn-router/conf/etc/service/sidecar/run
#!/bin/sh

[ ! -z ${ADDRS} ] && ARGS="-addrs ${ADDRS}"
[ ! -z ${SERVERS} ] && ARGS="${ARGS} -servers ${SERVERS}"

exec sidecar -name ${NAME} ${ARGS}
