# deploy/nfd/conf/etc/service/ndn-autoconfig-server/run
#!/bin/bash

if [ -z "$NDN_AUTOCONFIG_FACEURI" ]; then
    echo "NDN_AUTOCONFIG_FACEURI not set"
    exit 1
fi

exec ndn-autoconfig-server $NDN_AUTOCONFIG_PREFIXES $NDN_AUTOCONFIG_FACEURI
