# deploy/ndn-traffic/conf/service/ndn-traffic-server/run
#!/bin/sh

sleep 2
exec ndn-traffic-server /usr/local/etc/ndn/ndn-traffic-server.conf
